# DEBIAN/prerm
#!/bin/bash
if [ -d /run/systemd/system ]; then
    systemctl stop chonos-task.service
fi
systemctl disable chonos-task.service