# DEBIAN/prerm
#!/bin/bash
systemctl stop chonos-task.service
systemctl disable chonos-task.service
systemctl daemon-reload